# BTCOptionsTrading/backend/scripts/deployment/sentiment_api.service
[Unit]
Description=Sentiment Trading API Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/BTCOptionsTrading/backend
Environment="PATH=/root/BTCOptionsTrading/backend/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/root/BTCOptionsTrading/backend/venv/bin/python3 /root/BTCOptionsTrading/backend/sentiment_api.py
Restart=always
RestartSec=10
StandardOutput=append:/root/BTCOptionsTrading/backend/logs/sentiment_api.log
StandardError=append:/root/BTCOptionsTrading/backend/logs/sentiment_api.log

[Install]
WantedBy=multi-user.target
